# site/celery.conf
[program:celery]
command=/workspace/venv_site/bin/celery -A dmoj_celery worker
directory=/workspace/site
# You should create a dedicated user for celery to run under.
user=root
group=root
stdout_logfile=/tmp/celery.stdout.log
stderr_logfile=/tmp/celery.stderr.log
