# site/site.conf
[program:site]
command=/workspace/venv_site/bin/uwsgi --ini uwsgi.ini
directory=/workspace/site
stopsignal=QUIT
stdout_logfile=/tmp/site.stdout.log
stderr_logfile=/tmp/site.stderr.log
